# ready.service
[Unit]
Description=SCENARIO_IS_READY

[Service]
User=admin
ExecStart=echo SCENARIO_IS_READY
Restart=on-abort

[Install]
WantedBy=multi-user.target
