# ttyd.service
[Unit]
Description=TTYd remote access

[Service]
#User=admin
ExecStart=/usr/bin/ttyd login
Restart=on-abort

[Install]
WantedBy=multi-user.target
